# scripts/launchd/com.zerek.knowledge-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  R12.6.1 launchd-агент — авто-синк knowledge/ в GitHub раз в 30 минут.

  Установка (один раз):
    cp scripts/launchd/com.zerek.knowledge-sync.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.zerek.knowledge-sync.plist

  Удаление:
    launchctl unload ~/Library/LaunchAgents/com.zerek.knowledge-sync.plist
    rm ~/Library/LaunchAgents/com.zerek.knowledge-sync.plist

  Проверить что работает:
    launchctl list | grep zerek
    tail -20 /tmp/zerek_knowledge_sync.log
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zerek.knowledge-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/adil/Documents/ZEREK/scripts/sync_knowledge.sh</string>
    </array>

    <!-- Запускать каждые 30 минут (1800 секунд) -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Запустить сразу после load (полезно для проверки) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Если упал — не перезапускать сразу, ждать следующего интервала -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>WorkingDirectory</key>
    <string>/Users/adil/Documents/ZEREK</string>

    <!-- PATH должен включать /opt/homebrew/bin для gh (git credential helper) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/adil</string>
    </dict>

    <!-- stdout / stderr для дебага -->
    <key>StandardOutPath</key>
    <string>/tmp/zerek_knowledge_sync.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/zerek_knowledge_sync.stderr.log</string>
</dict>
</plist>
